# scripts/systemd/x3-desktop.service
[Unit]
Description=X3 Desktop Vite dev server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=10
WorkingDirectory=%h/Desktop/X3_ATOMIC_STAR/apps/x3-desktop
ExecStart=/usr/bin/env bash -lc 'npm run dev'
Environment=PORT=5173
Environment=CI=
StandardOutput=journal
StandardError=journal
KillMode=process

[Install]
WantedBy=default.target
