# .kilo/worktrees/luminous-pecorino/deployment/systemd/redis.service
[Unit]
Description=Redis In-Memory Data Store (for X3 Chain)
Documentation=man:redis-server(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=x3
Group=x3
ExecStart=/usr/bin/redis-server --port 6379
ExecStop=/usr/bin/redis-cli -p 6379 shutdown
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/redis /var/log/redis

[Install]
WantedBy=multi-user.target
