# .rc4-worktrees/old/deployment/systemd/x3-intelligence.service
[Unit]
Description=X3 Chain X3 Intelligence Dashboard
Documentation=man:x3-intelligence(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=x3
WorkingDirectory=/home/lojak/Desktop/x3-chain-master/apps/x3-intelligence
ExecStart=/usr/bin/npm run build && /usr/bin/npm run preview
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="PORT=5173"
Environment="HOST=0.0.0.0"
Environment="AUTH_SALT=x3-chain-production-2026"
Environment="VITE_API_URL=http://localhost:8000"
Environment="PATH=/home/x3/x3-chain/apps/x3-intelligence/node_modules/.bin:/usr/local/bin:/usr/bin:/bin"

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log/x3 /tmp

[Install]
WantedBy=multi-user.target
